# addon-AqaraPost/rootfs/etc/s6-overlay/s6-rc.d/init-nodered/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Node-RED
# Configures Node-RED before running
# ==============================================================================
declare port

if bashio::config.true 'renew_flow'; then
	bashio::log.info "renew_flow true"
    rm /config/flows.json && 	bashio::log.info "/config/flows.json deleted"|| bashio::log.info "/config/flows.json not deleted"

	else
	bashio::log.info "renew_flow false"
fi

if  bashio::fs.file_exists '/config/flows.json'; then
	bashio::log.info "......."
    bashio::log.info "flows.json already exist, the file will not be updated from config"
	bashio::log.info "......."
else
	#Install python for token generator aqara
	bashio::log.info "......."
	bashio::log.info "Install python3/pycryptodome for token generator Aqara ..."
	apk add python3 && apk add py3-pip && pip install requests && pip install pycryptodome && pip install requests-toolbelt  && bashio::log.info "Install python3/pycryptodome complete"|| bashio::log.info "install python3 error"
	
	#AqaraPost
	bashio::log.info "......."
	rm -rf /config/* && bashio::log.info "/config/flows.json not found, delete all folder and reset flow..." || bashio::log.info "......."
	bashio::log.info "......."
	
	bashio::log.info "Download json Aqara from GitHub"
	touch /etc/node-red/flows.json || bashio::log.info "...."
	bashio::log.info "......."
	
	FLOW_URL1="https://raw.githubusercontent.com/sdavides/AqaraPOST-Homeassistant/main/Aqara_G3_nodered.json"
	FLOW_NAME1="Aqara_G3_nodered"

	FLOW_URL2="https://raw.githubusercontent.com/sdavides/AqaraPOST-Homeassistant/main/Aqara_FP2_nodered.json"
	FLOW_NAME2="Aqara_FP2_nodered"

	FLOW_URL3="https://raw.githubusercontent.com/sdavides/AqaraPOST-Homeassistant/main/Aqara_G3_nodered_2device.json"
	FLOW_NAME3="Aqara_G3_nodered_2device"


	LUMI1=$(bashio::config 'lumi1Aqara') || bashio::log.info "error lumi1Aqara ......."
	LUMI2=$(bashio::config 'lumi1Aqara2') || bashio::log.info "error lumi1Aqara2 ......."
	LUMI3=$(bashio::config 'lumi1Aqara3') || bashio::log.info "error lumi1Aqara3 ......."
	
	LUMI1=$(echo "$LUMI1" | tr '[:upper:]' '[:lower:]' | sed -E 's/-//g; /^lumi1\./!s/^/lumi1./')
	LUMI2=$(echo "$LUMI2" | tr '[:upper:]' '[:lower:]' | sed -E 's/-//g; /^lumi1\./!s/^/lumi1./')
	LUMI3=$(echo "$LUMI3" | tr '[:upper:]' '[:lower:]' | sed -E 's/-//g; /^lumi1\./!s/^/lumi1./')
	
	TIMEZONE=$(bashio::config 'timezoneAqara') || bashio::log.info "error timezoneAqara ......."
	
    if [ "$LUMI1" = "" ]; then	
	bashio::log.warning "G3 lumi1 not set"
	else
	wget -q "$FLOW_URL1" --output-document=/etc/node-red/$FLOW_NAME1.json && bashio::log.info "Download AqaraPost $FLOW_NAME1.json complete..." || bashio::log.info "Error download AqaraPost $deviceType1 flow json"
	fi
	
    if [ "$LUMI2" = "" ]; then
	bashio::log.warning "FP2 lumi1 not set"
	else		
	wget -q "$FLOW_URL2" --output-document=/etc/node-red/$FLOW_NAME2.json && bashio::log.info "Download AqaraPost $FLOW_NAME2.json complete..." || bashio::log.info "Error download AqaraPost $deviceType2 flow json"
	fi

    if [ "$LUMI3" = "" ]; then	
	bashio::log.warning "G3_2 lumi1 not set"
	else
	wget -q "$FLOW_URL3" --output-document=/etc/node-red/$FLOW_NAME3.json && bashio::log.info "Download AqaraPost $FLOW_NAME3.json complete..." || bashio::log.info "Error download AqaraPost $deviceType3 flow json"
	fi
	
	ls /etc/node-red/Aqara_*.json || bashio::log.info "...."
		
	bashio::log.info "Starting AqaraPost-Homeassistant tokenGenerator ..."
	wget -q https://raw.githubusercontent.com/sdavides/AqaraPOST-Homeassistant/main/generatejson/AqaraPOST-tokenGenerator.py --output-document=/tmp/AqaraPOST-tokenGenerator.py  && chmod +x /tmp/AqaraPOST-tokenGenerator.py && ( echo $(bashio::config 'usernameAqara') ; echo $(bashio::config 'passwordAqara'); echo $(bashio::config 'serverAqara') )  | /tmp/AqaraPOST-tokenGenerator.py > /tmp/AqaraPOST-tokenGenerator.json || bashio::log.info "error......."
	cat /tmp/AqaraPOST-tokenGenerator.json | grep Username || bashio::log.info "Username Token error......."
	cat /tmp/AqaraPOST-tokenGenerator.json | grep Area || bashio::log.info "Area Token error......."

    if grep -q "Token:" /tmp/AqaraPOST-tokenGenerator.json; then
    TOKEN=$(grep "Token:" /tmp/AqaraPOST-tokenGenerator.json | sed 's/Token://')
    else
    bashio::log.warning "error Token..."
	TOKEN=
    fi

    if grep -q "AppID:" /tmp/AqaraPOST-tokenGenerator.json; then
    APPID=$(grep "AppID:" /tmp/AqaraPOST-tokenGenerator.json | sed 's/AppID://')
    else
    bashio::log.warning "error AppID..."
	APPID=
    fi

    if grep -q "Server:" /tmp/AqaraPOST-tokenGenerator.json; then
    SERVER=$(grep "Server:" /tmp/AqaraPOST-tokenGenerator.json | sed 's/Server://')
	cat /tmp/AqaraPOST-tokenGenerator.json | grep Server
    else
    bashio::log.warning "error Server..."
	SERVER=
    fi
	
	if  bashio::fs.file_exists '/etc/node-red/'$FLOW_NAME1'.json'; then
	sed -i "s/XXXXXXTOKENXXXXXXXXXXX/$TOKEN/g" /etc/node-red/$FLOW_NAME1.json && sed -i "s/XXXXXXAPPIDXXXXXXXXXXX/$APPID/g" /etc/node-red/$FLOW_NAME1.json  && sed -i "s/lumi1.XXXXXXXXXXXX/$LUMI1/g" /etc/node-red/$FLOW_NAME1.json && sed -i "s/rpc-ger.aqara.com/$SERVER/g" /etc/node-red/$FLOW_NAME1.json && sed -i "s/it-IT/$TIMEZONE/g" /etc/node-red/$FLOW_NAME1.json && bashio::log.info "tokenGenerator Insert value $FLOW_NAME1.json complete." || bashio::log.info "Error AqaraPost-Homeassistant tokenGenerator ..."
	fi
	
	if  bashio::fs.file_exists '/etc/node-red/'$FLOW_NAME2'.json'; then
	sed -i "s/XXXXXXTOKENXXXXXXXXXXX/$TOKEN/g" /etc/node-red/$FLOW_NAME2.json && sed -i "s/XXXXXXAPPIDXXXXXXXXXXX/$APPID/g" /etc/node-red/$FLOW_NAME2.json  && sed -i "s/lumi1.XXXXXXXXXXXX/$LUMI2/g" /etc/node-red/$FLOW_NAME2.json && sed -i "s/rpc-ger.aqara.com/$SERVER/g" /etc/node-red/$FLOW_NAME2.json && sed -i "s/it-IT/$TIMEZONE/g" /etc/node-red/$FLOW_NAME2.json && bashio::log.info "tokenGenerator Insert value $FLOW_NAME2.json complete." || bashio::log.info "Error AqaraPost-Homeassistant tokenGenerator ..."
	fi
	
	if  bashio::fs.file_exists '/etc/node-red/'$FLOW_NAME3'.json'; then
	sed -i "s/XXXXXXTOKENXXXXXXXXXXX/$TOKEN/g" /etc/node-red/$FLOW_NAME3.json && sed -i "s/XXXXXXAPPIDXXXXXXXXXXX/$APPID/g" /etc/node-red/$FLOW_NAME3.json  && sed -i "s/lumi1.XXXXXXXXXXXX/$LUMI3/g" /etc/node-red/$FLOW_NAME3.json && sed -i "s/rpc-ger.aqara.com/$SERVER/g" /etc/node-red/$FLOW_NAME3.json && sed -i "s/it-IT/$TIMEZONE/g" /etc/node-red/$FLOW_NAME3.json && bashio::log.info "tokenGenerator Insert value $FLOW_NAME3.json complete." || bashio::log.info "Error AqaraPost-Homeassistant tokenGenerator ..."
	fi
	
	jq -s 'add' /etc/node-red/Aqara_*json > /etc/node-red/flows.json && rm -rf /etc/node-red/Aqara_*json && bashio::log.info "flows.json created" || bashio::log.info "......." 
	rm -rf /tmp/AqaraPOST-tokenGenerator.py || bashio::log.info "......."
	rm -rf /tmp/AqaraPOST-tokenGenerator.json || bashio::log.info "......."
	
	bashio::log.info "......."
	bashio::log.info "End AqaraPost-Homeassistant flow" 
	bashio::log.info "......."

fi

	bashio::log.info "Starting Node-Red ..." 
	
# Migrate add-on data from the Home Assistant config folder,
# to the add-on configuration folder.
if ! bashio::fs.file_exists '/config/settings.js' \
    && bashio::fs.file_exists '/homeassistant/node-red/settings.js'; then
    shopt -s dotglob
    mv /homeassistant/node-red/* /config/ \
        || bashio::exit.nok "Failed to migrate Node-RED configuration"
fi

# Ensure configuration exists
if ! bashio::fs.file_exists '/config/settings.js'; then
    mkdir -p /config/nodes \
        || bashio::exit.nok "Failed to create node-red configuration directory"

    # Copy in template files
    cp /etc/node-red/flows.json /config/
    cp /etc/node-red/settings.js /config/

    # Create random flow id
    id=$(node -e "console.log((2+Math.random()*4294967295).toString(26));")
    sed -i "s/%%ID%%/${id}/" "/config/flows.json"
fi

# Pass in port & SSL settings
port=$(bashio::addon.port 80)
sed -i "s/%%PORT%%/${port:-80}/" "/opt/node_modules/node-red-dashboard/nodes/ui_base.html"
if ! bashio::var.has_value "${port}"; then
    bashio::log.warning
    bashio::log.warning "Direct access mode is disabled, Node-RED Dashboard"
    bashio::log.warning "will not work!"
    bashio::log.warning
    bashio::log.warning "Please assign a port in the Network section of this"
    bashio::log.warning "add-on configuration."
    bashio::log.warning
fi

if bashio::config.true 'ssl'; then
    sed -i "s/%%SSL%%/true/" "/opt/node_modules/node-red-dashboard/nodes/ui_base.html"
else
    sed -i "s/%%SSL%%/false/" "/opt/node_modules/node-red-dashboard/nodes/ui_base.html"
fi

# Ensures conflicting Node-RED packages are absent
cd /config || bashio::exit.nok "Could not change directory to Node-RED"
if bashio::fs.file_exists "/config/package.json"; then
    npm uninstall \
        node-red-contrib-home-assistant \
        node-red-contrib-home-assistant-llat \
        node-red-contrib-home-assistant-ws \
            || bashio::log.warning "Failed un-installing conflicting packages"
fi
